# etc/pi/update-motd.d/20-raspberry-pi
#!/bin/sh
printf "
        ⟨  Raspberry  ⟩        $(uname -srm)
      ________   ________ 
     /        \\ /        \\
    /         /_/       / 
   //      __//         / 
   \\\\\\\\_____/   \\\\\\\\_______/  

"
